# HAPCStat.timer
[Timer]
OnBootSec=2minutes
OnUnitActiveSec=2minutes
Unit=HAPCStat.service